# reffsd.service.in
# SPDX-FileCopyrightText: 2025 Tom Haynes <loghyr@gmail.com>
# SPDX-License-Identifier: AGPL-3.0-or-later

[Unit]
Description=Reffs NFS Daemon
Wants=network-online.target
Requires=rpcbind.socket
After=network-online.target rpcbind.socket

[Service]
Type=simple
Environment="LIBC_FATAL_STDERR_=1"
Environment="MALLOC_CONF=prof:false"
Environment="ASAN_OPTIONS=quarantine_size_mb=256:detect_leaks=0:halt_on_error=0:handle_abort=1"
Environment="UBSAN_OPTIONS=halt_on_error=0:print_stacktrace=1"
ExecStart=@FULL_BINDIR@/reffsd --config=@FULL_SYSCONFDIR@/reffs/reffsd.toml
Restart=on-failure
RestartSec=5
LimitCORE=infinity
LimitNOFILE=65536
StateDirectory=reffs
LogsDirectory=reffs
SyslogIdentifier=reffsd

[Install]
WantedBy=multi-user.target
